# systemd/audit-webhook.service
[Unit]
Description=Audit Webhook Container
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
Environment="IMAGE=ghcr.io/OWNER/REPO:main"
EnvironmentFile=/srv/audit-webhook/.env
ExecStartPre=-/usr/bin/docker pull ${IMAGE}
ExecStart=/usr/bin/docker run --rm \
    --name audit-webhook \
    --env-file /srv/audit-webhook/.env \
    -p 8080:8080 \
    ${IMAGE}
ExecStop=/usr/bin/docker stop audit-webhook
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
